# treehole-email-bot.service
[Unit]
Description=PKU Treehole Email Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/pku-treehole-search-agent
ExecStart=/usr/bin/python3 /home/ubuntu/pku-treehole-search-agent/bot_email.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/pku-treehole-search-agent/logs/bot.log
StandardError=append:/home/ubuntu/pku-treehole-search-agent/logs/bot.error.log

[Install]
WantedBy=multi-user.target
